# roles/base/files/login-block.service
[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/usr/local/bin/login-block

StandardOutput=journal+console
