# roles/base/files/reset-ent-sandbox.service
[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/usr/local/bin/reset-ent-sandbox

StandardOutput=journal+console
